# mpd/.config/systemd/user/mpd.service
[Unit]
Description=Some daemon for music

[Service]
Type=forking
ExecStart=/usr/bin/mpd

[Install]
WantedBy=default.target
